# scripts/CcspInangoplugComponent.path
[Unit]
Description=Checks that the Psm initialized

[Path]
PathExists=/tmp/psm_initialized

[Install]
WantedBy=multi-user.target
