# scripts/check_inangoplug_enabled.path
[Unit]
Description=OvS enabled checking
After=logagent.service

[Path]
PathExists=/tmp/utopia_inited

[Install]
WantedBy=multi-user.target
