# scripts/connect_inangoplug.path
[Unit]
Description=Checks that the PandM initialized

[Path]
PathExists=/tmp/pam_initialized

[Install]
WantedBy=multi-user.target

